# install/etc/profile.d/pvzdweb.sh
source /opt/venv/satosa_rpmgr/bin/activate
