# install/etc/profile.d/satosa_rpmgr.sh
#!/usr/bin/env bash
source /opt/venv/satosa_rpmgr/bin/activate
